# app/config/_os/usr/lib/systemd/system/ashuffle.service
[Unit]
Description=ashuffle for mpd
After=mpd.service network.target

[Service]
ExecStartPre=/usr/bin/sleep 2
ExecStart=/usr/bin/ashuffle
ExecReload=/usr/bin/kill -HUP $MAINPID
KillMode=process
Restart=always
RestartSec=10
TimeoutStartSec=20
StartLimitInterval=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
