# app/config/_os/usr/lib/systemd/system/redis.socket
[Socket]
ListenStream=/run/redis/socket
ListenStream=6379
Backlog=5
KeepAlive=true
PassCredentials=true

[Install]
WantedBy=sockets.target
