# app/config/defaults/rune_SSM_wrk.service
[Unit]
Description=RuneAudio Playback Worker 
After=network.target redis.target rune_PL_wrk

[Service]
ExecStart=/var/www/command/rune_SSM_wrk
TimeoutSec=0
Restart=always
RestartSec=1
StartLimitInterval=30
StartLimitBurst=20

[Install]
WantedBy=multi-user.target
