# app/config/defaults/etc/udev/rules.d/60-bluetooth.rules
ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="/srv/http/command/bt-info.sh start"
ACTION=="remove", SUBSYSTEM=="bluetooth", RUN+="/srv/http/command/bt-info.sh stop"
