# app/config/defaults/etc/udev/rules.d/99-a2dp.rules
KERNEL=="input[0-9]*", RUN+="/srv/http/command/a2dp-autoconnect"
